# systemd/nvd_repo_rebuild.service.template
[Unit]
Description=NVD-Repo Rebuild

[Service]
Type=exec
User=<USER>
Group=<USER>
WorkingDirectory=<ROOT OF THIS PROJECT>
ExecStart=<ROOT OF THIS PROJECT>/scripts/systemd_timer_rebuild.bash
